# saruman_of_many_colorsbot/bot.service
[Unit]
Description=bot
After=default.target
[Service]
User=admin
ExecStart=/home/SANITIZED/Documents/tmux-bot.sh
[Install]
WantedBy=default.target